# opscache/etc/supervisor/conf.d/squid.conf
[program:squid]
command=sh /start_squid
